# .config/systemd/user/ipns-publish.timer
[Unit]
Description=IPNS publish timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=3h

[Install]
WantedBy=timers.target
