# .config/systemd/user/ipns-publish.service
[Unit]
Description=Publish latest CID to IPNS

[Service]
ExecStart=/home/ndavd/bin/ipns-publish
